# debian/isc-dhcp-relay.prerm
#!/bin/sh -e
#
# $Id: dhcp3-relay.prerm,v 1.1.1.1 2002/05/21 00:07:44 peloy Exp $
#

if ( [ "$1" = "upgrade" ] || [ "$1" = "remove" ] ) && \
	[ -L /usr/doc/dhcp3-relay ]; then
	rm -f /usr/doc/dhcp3-relay
fi

invoke-rc.d dhcp3-relay stop

#DEBHELPER#

exit 0
